# templates/kibana.upstart.erb
description 'Kibana service'

start on static-network-up

exec <%= @install_path %>/kibana/bin/kibana
